# roles/pg_backup/files/pg_backup.timer
[Unit]
Description=Backup PostgreSQL databases daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
