# roles/gimp_rsyncd/files/rsyncd@.service
[Unit]
Description=Rsync daemon
After=network.target

[Service]
ExecStart=/usr/bin/rsync --daemon --config=/etc/rsyncd.conf
StandardInput=socket