# roles/gimp_rsyncd/files/rsyncd.socket
[Unit]
Description=rsyncd Service Sockets

[Socket]
ListenStream=873
Accept=yes

[Install]
WantedBy=sockets.target