# roles/gitlab/templates/gitlab-mirror2github.service
[Unit]
Description=Mirror repositories to GitLab

[Service]
Type=oneshot
User=sysadmin
ExecStart=/home/admin/bin/gitlab/mirror2github
