# roles/nsd/templates/automatic-dnssec-signzone.service.j2
[Unit]
Description=Automatic DNSSEC signzone

[Service]
Type=oneshot
User=root
WorkingDirectory=/tmp
ExecStart=/usr/local/bin/automatic-dnssec-signzone