# roles/nsd/templates/automatic-dnssec-signzone.timer.j2
[Unit]
Description=Automatic DNSSEC signzone

[Timer]
OnCalendar=weekly
RandomizedDelaySec=12hours
Persistent=true

[Install]
WantedBy=timers.target
